# host-mcp-gateway/launchd/brain.host-mcp-gateway.plist.sample
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>brain.host-mcp-gateway</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/host-mcp-gateway</string>
      <string>-config</string>
      <string>/Users/CHANGEME/.config/brain/host-mcp-gateway.json</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>OTEL_EXPORTER_OTLP_ENDPOINT</key>
      <string>http://localhost:4317</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/CHANGEME/Library/Logs/brain-host-mcp-gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/CHANGEME/Library/Logs/brain-host-mcp-gateway.error.log</string>
  </dict>
</plist>
